# roles/nagios/files/cron.d/nagios_conf
*/30 * * * * nagios /home/nagios/githome/nagios_conf/update > /dev/null 2>&1
