# roles/__archived/deluge/templates/deluge-web.service.j2
[Unit]
Description=Deluge Web
After=network.target

[Service]
User=deluge
Group={{ global_rw_group }}
ExecStart=/usr/local/bin/deluge-web -d -L error
Restart=on-failure
SuccessExitStatus=3 4
RestartForceExitStatus=3 4
UMask=0027

[Install]
WantedBy=multi-user.target

